# experiment/tools/supervisor/k8sdashboard.conf
[program:k8s-dashboard-8282]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8282
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8282.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8282.log
stderr_logfile_maxbytes=100MB

[program:k8s-dashboard-8283]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8283
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8283.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8283.log
stderr_logfile_maxbytes=100MB

[program:k8s-dashboard-8284]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8284
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8284.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8284.log
stderr_logfile_maxbytes=100MB

[program:k8s-dashboard-8285]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8285
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8285.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8285.log
stderr_logfile_maxbytes=100MB

[program:k8s-dashboard-8286]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8286
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8286.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8286.log
stderr_logfile_maxbytes=100MB

[program:k8s-dashboard-8287]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8287
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8287.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8287.log
stderr_logfile_maxbytes=100MB

[program:k8s-dashboard-8288]
command=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/main-8288
autorestart=true
redirect_stderr=true
stdout_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8288.log
stdout_logfile_maxbytes=100MB
stderr_logfile=/srv/go/src/gitlab.luojilab.com/k8s-dashboard/logs/k8s-dashboard-8288.log
stderr_logfile_maxbytes=100MB