# experiment/tools/supervisor/phpcs.conf
[program:phpcsServer]
command = /usr/bin/php /data/wwwroot/test.php
autostart = true
startsecs = 5
autorestart = true
startretries = 3
user = root
redirect_stderr = true
stdout_logfile_maxbytes = 20MB
stdout_logfile_backups = 10
stdout_logfile = /var/log/supervisor/phpcsServer.log
stopasgroup=true